# docker_images/php_worker/app/crontab
#!/bin/sh
* * * * * root bash /app/rnd.sh
#Line for correct work of cron-record
